# recipes-spektrum/spektrum-sbc-oled/files/spektrum-oled-status.service
[Unit]
Description=Spektrum OLED status display
PartOf=spektrum-runtime.target
After=spektrum-first-boot.service network.target
Wants=spektrum-first-boot.service

[Service]
Type=simple
EnvironmentFile=-/etc/spektrum/spektrum.env
ExecStart=/usr/bin/bash -lc '/usr/bin/python3 ${SPEKTRUM_ROOT:-/opt/spektrum}/scripts/sbc/oled_status.py --state-db ${SPEKTRUM_STATE_DB:-/var/lib/spektrum/device_state.db} --i2c-port ${SPEKTRUM_I2C_PORT:-0} --i2c-address ${SPEKTRUM_I2C_ADDRESS:-0x3C} --width ${SPEKTRUM_OLED_WIDTH:-128} --height ${SPEKTRUM_OLED_HEIGHT:-32} --interval ${SPEKTRUM_OLED_INTERVAL:-2} --page-seconds ${SPEKTRUM_OLED_PAGE_SECONDS:-4} --rotate ${SPEKTRUM_OLED_ROTATE:-0}'
Restart=always
RestartSec=2

[Install]
WantedBy=spektrum-runtime.target
